# restart_run.timer
[Unit]
Description=restart scan broker
[Timer]
OnUnitActiveSec=3600s
OnBootSec=0s

[Install]
WantedBy=scan.service
